# arxiv-daily-push/deploy/cloudflare/launchd/com.linze.adp.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- Cloudflare Tunnel 常驻 connector（远程管理式，隧道名 adp）。
     令牌在 ~/.cloudflared/adp-tunnel-token（0600，不入库）；
     安装：cp 到 ~/Library/LaunchAgents/ 后 launchctl bootstrap gui/$UID <plist> -->
<dict>
  <key>Label</key><string>com.linze.adp.tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/linzezhang/Documents/Codex/main_worktree/CodexProject/adp/arxiv-daily-push/var/bin/cloudflared</string>
    <string>tunnel</string>
    <string>--no-autoupdate</string>
    <string>run</string>
    <string>--token-file</string>
    <string>/Users/linzezhang/.cloudflared/adp-tunnel-token</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key>
  <string>/Users/linzezhang/Documents/Codex/main_worktree/CodexProject/adp/arxiv-daily-push/var/log/adp-tunnel.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/linzezhang/Documents/Codex/main_worktree/CodexProject/adp/arxiv-daily-push/var/log/adp-tunnel.log</string>
</dict>
</plist>
