# Alpha/deploy/systemd/alpha-trading-worker.service
[Unit]
Description=Alpha 交易 Worker(SIMULATE/受门禁 REAL;失败关闭)
After=alpha-opend.service postgresql.service
Requires=alpha-opend.service

[Service]
User=alpha
WorkingDirectory=/opt/alpha/MetaDatabase/Alpha
EnvironmentFile=/opt/alpha/env
ExecStart=/opt/alpha/venv/bin/python -m backend.app.workers.main_trading
# 看门狗:装配完成报 READY,每拍喂 WATCHDOG;5 分钟没喂=卡死=杀掉重启
Type=notify
NotifyAccess=main
WatchdogSec=300
TimeoutStartSec=300
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
