# Alpha/deploy/systemd/alpha-alert@.service
# 定时任务失败自告警(OnFailure 钩子;%i = 失败的单元名)
# Type=oneshot;无任何激活权限,只入队一封人话告警邮件。
[Unit]
Description=Alpha 单元失败自告警(%i)

[Service]
Type=oneshot
User=alpha
Group=alpha
WorkingDirectory=/opt/alpha/MetaDatabase/Alpha
EnvironmentFile=/opt/alpha/env
ExecStart=/opt/alpha/venv/bin/python scripts/notify_unit_failure.py %i
TimeoutStartSec=120
